# rootfs/service/rtorrent/finish
#!/bin/sh

echo
echo "rTorrent exited."
