# rootfs/service/rtorrent/run
#!/bin/bash

# Remove the lock, otherwise rTorrent will get stuck in a loop and won't start.
rm -f /rtorrent/.session/rtorrent.lock

# Turn to rTorrent
exec rtorrent
